# docker/supervisor-mongo.conf
[program:mongo]
command=mongod --dbpath=/data/db
user=root